# systemd/dfan-tray.service
[Unit]
Description=dfan system tray icon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/sbin/dfan-tray
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
